# netxms-reporting.service
[Unit]
Description=NetXMS reporting server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nxreportd
Restart=always
UMask=0077

[Install]
WantedBy=multi-user.target
Alias=nxreportd.service
